# cronjob-tmpwatch.service
[Unit]
Description=tmpwatch cronjob

[Service]
ExecStart=/usr/lib/tmpwatch
