# cronjob-tmpwatch.timer
[Unit]
Description=tmpwatch cronjob timer

[Timer]
OnCalendar=hourly

[Install]
WantedBy=cronjobs.target
