# tmpwatch.cron
#!/bin/sh
#########
# Some defaults:
OPTIONS="-fq"
if [ -f /etc/sysconfig/tmpwatch ]; then
	. /etc/sysconfig/tmpwatch
fi

########
for file in /etc/tmpwatch/*.conf; do
	grep ^/ ${file} | while read DIRECTORY TIME OPTS; do
		OPTIONS="${OPTIONS} ${OPTS}"
		/usr/sbin/tmpwatch ${OPTIONS} ${TIME} ${DIRECTORY}
	done
done

