# tmpwatch.crontab
MAILTO=root

# tmpwatch to run hourly
02 * * * * root /usr/lib/tmpwatch
