# tmpwatch.crontab
MAILTO=root

# tmpwatch to run daily
02 5 * * * root /usr/lib/tmpwatch
